# .config/profile.d/20_antidoton.sh
test -f "${XDG_DATA_DIR:-$HOME/.local/share}/antidoton/antidoton.sh" && . "${XDG_DATA_DIR:-$HOME/.local/share}/antidoton/antidoton.sh"