# .config/profile.d/20_asdf.sh
#!/bin/sh
test -f "${XDG_DATA_HOME:-$HOME/.local/share}/_asdf/asdf.sh" && . "${XDG_DATA_HOME:-$HOME/.local/share}/_asdf/asdf.sh"